# og-helper.service
[Unit]
Description=OpenGradient Helper Telegram Bot
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/og-helper-bot
ExecStart=/root/og-helper-bot/venv/bin/python bot.py
EnvironmentFile=/root/og-helper-bot/.env

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=og-helper

[Install]
WantedBy=multi-user.target
